# services/upload-benchmarks.service
[Unit]
Description=Polars Remote Benchmark Uploader
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=5

[Service]
Type=oneshot
EnvironmentFile=/etc/systemd/system/upload-benchmarks.env
WorkingDirectory=/path/to/polars-remove-benchmark
ExecStart=/path/to/polars-remove-benchmark/upload-graphs.sh
User=polars
Group=polars

[Install]
WantedBy=multi-user.target