# services/upload-benchmarks.path
[Path]
Unit=upload-benchmarks.service
PathChanged=/path/to/polars-remote-benchmark/upload-probe

[Install]
WantedBy=multi-user.target
